# supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/usr/bin/redis-server /etc/redis/redis.conf
